# deploy/infinitee.service
# /etc/systemd/system/infinitee.service
# Enable: sudo systemctl enable infinitee
# Start:  sudo systemctl start infinitee
# Status: sudo systemctl status infinitee
# Logs:   sudo journalctl -u infinitee -f

[Unit]
Description=Infinitee Xclusives — Gunicorn Application Server
After=network.target mysql.service
Requires=mysql.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/path/to/infinitee_app
EnvironmentFile=/path/to/infinitee_app/.env
ExecStart=/path/to/venv/bin/gunicorn \
    --workers 4 \
    --worker-class sync \
    --worker-connections 1000 \
    --timeout 120 \
    --keep-alive 5 \
    --bind 127.0.0.1:8000 \
    --access-logfile /path/to/infinitee_app/logs/gunicorn-access.log \
    --error-logfile  /path/to/infinitee_app/logs/gunicorn-error.log \
    --log-level info \
    infinitee.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
